# deploy/systemd/citadel-watchdog.service
[Unit]
Description=Citadel Quantum Trader – Watchdog & Automatic Fail‑Over
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/opt/cqt/watchdog.sh
Restart=always
RestartSec=5
User=root
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
