# etc/systemd/system/dr-drill.timer
[Unit]
Description=Run Citadel DR‑drill every night at 02:00 UTC

[Timer]
OnCalendar=*-*-* 02:00:00 UTC
Persistent=true
Unit=dr-drill.service

[Install]
WantedBy=timers.target
